# assembly/src/main/deb/DEBIAN/prerm
#!/bin/sh
# prerm script for daisy-pipeline2
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <prerm> `remove'
#        * <old-prerm> `upgrade' <new-version>
#        * <new-prerm> `failed-upgrade' <old-version>
#        * <conflictor's-prerm> `remove' `in-favour' <package> <new-version>
#        * <deconfigured's-prerm> `deconfigure' `in-favour'
#          <package-being-installed> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package

SYSTEM_PROPERTIES=${deb.pipeline.home}/etc/pipeline.properties
CLI_CONFIG_YML=${deb.pipeline.home}/cli/config.yml

case "$1" in
    
    upgrade|remove)
        
        # -----------------
        # upgrade or remove
        # -----------------
        
        # mark configuration files (will be restored in postrm)
        
        if [ "$1" = remove ]; then
            [ -e $SYSTEM_PROPERTIES ] && cp $SYSTEM_PROPERTIES $SYSTEM_PROPERTIES.keep
            [ -e $CLI_CONFIG_YML ] && cp $CLI_CONFIG_YML $CLI_CONFIG_YML.keep
        fi
        
        # stop the daemon
        /etc/init.d/daisy-pipeline2 stop
        
        # remove the System V style init script links
        update-rc.d -f daisy-pipeline2 remove >/dev/null
        
    ;;
    
    deconfigure|failed-upgrade)
    ;;

    *)
        echo "prerm called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
